# lucee-nginx/supervisord.conf
[supervisord]
nodaemon=true

[program:tomcat]
command=/usr/local/tomcat/bin/catalina.sh run

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"